# ops/systemd/grafana-proxy.service
[Unit]
Description=Expose Grafana from Minikube on VM port 3000
After=network-online.target docker.service
Wants=network-online.target

[Service]
Type=simple
User=vagrant
WorkingDirectory=/vagrant/student-management
ExecStart=/vagrant/student-management/ops/minikube-nodeport-proxy.sh 3000 32000
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
