# ops/systemd/prometheus-proxy.service
[Unit]
Description=Expose Prometheus from Minikube on VM port 9090
After=network-online.target docker.service
Wants=network-online.target

[Service]
Type=simple
User=vagrant
WorkingDirectory=/vagrant/student-management
ExecStart=/vagrant/student-management/ops/minikube-nodeport-proxy.sh 9090 30090
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
